# rootfs/etc/services.d/cbbLocalManagement/run
/opt/local/CloudBerry Backup/bin/cbbLocalManagement