# rootfs/etc/services.d/cbbLocalManagement/run
/opt/local/MSP360 Backup/bin/cbbLocalManagement